# etc/inittab
#
# A sample of the new, SYSV compatible inittab.
#
# (This one runs on my machine every day)
#
# Level to run in. Set to 5 or 6 to allow serial port logins.
# If you comment this out, the system will ask you for a runlevel
# when it is booted.
#
id:6:initdefault:
#
# Sysinit: takes place only once, right after system boot, *before*
# possibly going single-user.
#
si::sysinit:/etc/update &
#
# boot & bootwait take place once too, but *after* possibly going
# single user.
#
rc::bootwait:/etc/rc
#
# Run the cron deamon.
#
#cr::boot:/usr/bin/crond
#
# Normal levels: 1-4 = virtual consoles 5,6 = COM ports.
# Ofcourse you can change this to suit your taste.
#
# BEWARE: where is your getty? in /bin or in /etc?
#
c1:123456:respawn:/bin/getty 9600 tty1
c2:23456:respawn:/bin/getty 9600 tty2
c3:3456:respawn:/bin/getty 9600 tty3
c4:456:respawn:/bin/getty 9600 tty4
c5:456:respawn:/bin/getty 9600 tty5
c6:456:respawn:/bin/getty 9600 tty6
#c7:6:respawn:/bin/getty 9600 ttys1
#c8:6:respawn:/bin/uugetty 9600 ttys0
#c6:6:respawn:/bin/getty 19200 ttys1
#
# The powerfail routines, that take care of a graceful shutdown.
# In our case, CTRL-ALT-DEL was pressed.
#
pf::powerwait:/etc/shutdown -rf now > /dev/console < /dev/console 2>&1
