# etc/rc.local
/usr/bin/setterm -blank 10
if [ -x /etc/rc.net ]; then 
	/bin/hostname `cat /etc/host`
	/etc/rc.net 
	(sleep 2;          # wait for daemons to get a foothold 
	mount -av -t nfs)&
else
	/bin/hostname softland
fi 

if [ -x /etc/clock ]; then
	/etc/clock -s
fi

if [ -x /etc/syslogd ]; then
	syslogd  # turn on logging of su's, logins etc.
fi

if [ -x /etc/syslogk ]; then
	syslogk on > /usr/log/kernel&  # log all kernel messages.
fi
